# services/calcurse/config/calcurse.service
[Unit]
Description=Run calcurse-caldav to sync calcurse

[Service]
Type=oneshot
ExecStart=/usr/bin/calcurse-caldav
